# systemd/emotion-backend.service
[Unit]
Description=Emotion Companion Remote Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
SupplementaryGroups=audio video i2c gpio
EnvironmentFile=-/etc/default/emotion-pi
Environment=PYTHONUNBUFFERED=1
Environment=OMP_NUM_THREADS=1
Environment=OPENBLAS_NUM_THREADS=1
Environment=NUMEXPR_NUM_THREADS=1
WorkingDirectory=__PROJECT_ROOT__
ExecStart=/bin/sh -lc 'exec "__PROJECT_ROOT__/.venv/bin/python" server_backend/run_server.py'
Restart=always
RestartSec=2
TasksMax=128
Nice=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
